# HarfangLab/harfanglab/tests/alert_1.json
{
  "input": {
    "message": "{\n  \"level_int\": 30,\n  \"detection_origin\": \"agent\",\n  \"threat_key\": \"2971\",\n  \"log_type\": \"alert\",\n  \"rule_name\": \"File Added/Modified in Startup Directory\",\n  \"status\": \"new\",\n  \"aggregation_key\": \"e9cb2440a0ef88bdf8133cb7d5e41b66a903ba84df49659a05720ee1ac1caa33\",\n  \"tags\": [\n    \"attack.persistence\",\n    \"attack.t1547.001\"\n  ],\n  \"tenant\": \"\",\n  \"level\": \"medium\",\n  \"rule_id\": \"73d2962a-48ed-4fc1-937c-0a5bc688a072\",\n  \"quarantine\": 4,\n  \"alert_subtype\": \"process\",\n  \"threat_values\": [],\n  \"details_file\": {\n    \"target_filename\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Roaming\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\DeepL auto-start.lnk\"\n  },\n  \"threat_type\": \"new\",\n  \"groups\": [],\n  \"alert_type\": \"sigma\",\n  \"process\": {\n    \"pe_info\": {\n      \"file_description\": \"Windows Explorer\",\n      \"product_version\": \"10.0.19041.1706\",\n      \"file_version\": \"10.0.19041.1706 (WinBuild.160101.0800)\",\n      \"company_name\": \"Microsoft Corporation\",\n      \"product_name\": \"Microsoft\u00ae Windows\u00ae Operating System\",\n      \"original_filename\": \"EXPLORER.EXE\",\n      \"internal_name\": \"explorer\",\n      \"pe_timestamp\": \"2019-11-01T20:27:45.000Z\",\n      \"legal_copyright\": \"\u00a9 Microsoft Corporation. All rights reserved.\"\n    },\n    \"logonid\": 1111111,\n    \"commandline\": \"C:\\\\Windows\\\\Explorer.EXE\",\n    \"parent_unique_id\": \"27717f1b-edb4-49d7-aa3b-aab2e0beb389\",\n    \"hashes\": {\n      \"sha256\": \"01ba4719c80b6fe911b091a7c05124b64eeece964e09c058ef8f9805daca546b\",\n      \"sha1\": \"adc83b19e793491b1c6ea0fd8b46cd9f32e592fc\",\n      \"md5\": \"68b329da9893e34099c7d8ad5cb9c940\"\n    },\n    \"log_type\": \"process\",\n    \"usersid\": \"S-1-5-21-1111111111-2222222222-3333333333-44444\",\n    \"status\": 0,\n    \"parent_integrity_level\": \"Medium\",\n    \"size\": 5114880,\n    \"pe_timestamp_int\": 1572640065,\n    \"parent_image\": \"C:\\\\Windows\\\\System32\\\\userinit.exe\",\n    \"process_unique_id\": \"762b5b06-dfcf-431e-8352-7b83d6ff7ed6\",\n    \"session\": 1,\n    \"ioc_matches\": [],\n    \"pe_timestamp\": \"2019-11-01T20:27:45.000Z\",\n    \"grandparent_image\": \"C:\\\\Windows\\\\System32\\\\winlogon.exe\",\n    \"fake_parent_image\": \"\",\n    \"signed\": true,\n    \"process_name\": \"explorer.exe\",\n    \"sigma_rule_content\": \"\",\n    \"pid\": 11560,\n    \"fake_ppid\": 0,\n    \"fake_parent_commandline\": \"\",\n    \"grandparent_integrity_level\": \"System\",\n    \"current_directory\": \"C:\\\\Windows\\\\system32\\\\\",\n    \"image_name\": \"C:\\\\Windows\\\\explorer.exe\",\n    \"create_time\": \"2024-01-15T08:12:43.045Z\",\n    \"ppid\": 6452,\n    \"grandparent_commandline\": \"winlogon.exe\",\n    \"log_platform_flag\": 0,\n    \"parent_commandline\": \"C:\\\\Windows\\\\system32\\\\userinit.exe\",\n    \"ancestors\": \"C:\\\\Windows\\\\System32\\\\userinit.exe|C:\\\\Windows\\\\System32\\\\winlogon.exe\",\n    \"integrity_level\": \"Medium\",\n    \"username\": \"EXAMPLE\\\\jdoe\"\n  },\n  \"alert_unique_id\": \"44c633d9-b38d-4acb-87a5-7db9bd8ab38a\",\n  \"type\": \"rtlogs\",\n  \"execution\": 0,\n  \"alert_time\": \"2024-01-15T08:13:47.621+00:00\",\n  \"image_name\": \"C:\\\\Windows\\\\explorer.exe\",\n  \"mitre_cells\": [\n    \"persistence__t1547.001\"\n  ],\n  \"@event_create_date\": \"2024-01-15T08:13:47.621Z\",\n  \"@version\": \"1\",\n  \"maturity\": \"stable\",\n  \"agent\": {\n    \"additional_info\": null,\n    \"agentid\": \"d82e1ff6-f268-42a1-9091-df230fb3c85e\",\n    \"groups\": [\n      {\n        \"id\": \"ba5a8596-5bfc-401a-a46e-1111f7d35a23\",\n        \"name\": \"EXAMPLE\"\n      }\n    ],\n    \"osversion\": \"10.0.19042\",\n    \"osproducttype\": \"Windows 10 Pro\",\n    \"domainname\": \"EXAMPLE\",\n    \"distroid\": null,\n    \"ostype\": \"windows\",\n    \"hostname\": \"PL-3049\",\n    \"version\": \"3.2.5\",\n    \"dnsdomainname\": \"example.org\",\n    \"domain\": null\n  },\n  \"@timestamp\": \"2024-01-15T08:14:59.497966Z\",\n  \"msg\": \"Detects when a file is added or modified in the startup directory\"\n}"
  },
  "expected": {
    "message": "{\n  \"level_int\": 30,\n  \"detection_origin\": \"agent\",\n  \"threat_key\": \"2971\",\n  \"log_type\": \"alert\",\n  \"rule_name\": \"File Added/Modified in Startup Directory\",\n  \"status\": \"new\",\n  \"aggregation_key\": \"e9cb2440a0ef88bdf8133cb7d5e41b66a903ba84df49659a05720ee1ac1caa33\",\n  \"tags\": [\n    \"attack.persistence\",\n    \"attack.t1547.001\"\n  ],\n  \"tenant\": \"\",\n  \"level\": \"medium\",\n  \"rule_id\": \"73d2962a-48ed-4fc1-937c-0a5bc688a072\",\n  \"quarantine\": 4,\n  \"alert_subtype\": \"process\",\n  \"threat_values\": [],\n  \"details_file\": {\n    \"target_filename\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Roaming\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\DeepL auto-start.lnk\"\n  },\n  \"threat_type\": \"new\",\n  \"groups\": [],\n  \"alert_type\": \"sigma\",\n  \"process\": {\n    \"pe_info\": {\n      \"file_description\": \"Windows Explorer\",\n      \"product_version\": \"10.0.19041.1706\",\n      \"file_version\": \"10.0.19041.1706 (WinBuild.160101.0800)\",\n      \"company_name\": \"Microsoft Corporation\",\n      \"product_name\": \"Microsoft\u00ae Windows\u00ae Operating System\",\n      \"original_filename\": \"EXPLORER.EXE\",\n      \"internal_name\": \"explorer\",\n      \"pe_timestamp\": \"2019-11-01T20:27:45.000Z\",\n      \"legal_copyright\": \"\u00a9 Microsoft Corporation. All rights reserved.\"\n    },\n    \"logonid\": 1111111,\n    \"commandline\": \"C:\\\\Windows\\\\Explorer.EXE\",\n    \"parent_unique_id\": \"27717f1b-edb4-49d7-aa3b-aab2e0beb389\",\n    \"hashes\": {\n      \"sha256\": \"01ba4719c80b6fe911b091a7c05124b64eeece964e09c058ef8f9805daca546b\",\n      \"sha1\": \"adc83b19e793491b1c6ea0fd8b46cd9f32e592fc\",\n      \"md5\": \"68b329da9893e34099c7d8ad5cb9c940\"\n    },\n    \"log_type\": \"process\",\n    \"usersid\": \"S-1-5-21-1111111111-2222222222-3333333333-44444\",\n    \"status\": 0,\n    \"parent_integrity_level\": \"Medium\",\n    \"size\": 5114880,\n    \"pe_timestamp_int\": 1572640065,\n    \"parent_image\": \"C:\\\\Windows\\\\System32\\\\userinit.exe\",\n    \"process_unique_id\": \"762b5b06-dfcf-431e-8352-7b83d6ff7ed6\",\n    \"session\": 1,\n    \"ioc_matches\": [],\n    \"pe_timestamp\": \"2019-11-01T20:27:45.000Z\",\n    \"grandparent_image\": \"C:\\\\Windows\\\\System32\\\\winlogon.exe\",\n    \"fake_parent_image\": \"\",\n    \"signed\": true,\n    \"process_name\": \"explorer.exe\",\n    \"sigma_rule_content\": \"\",\n    \"pid\": 11560,\n    \"fake_ppid\": 0,\n    \"fake_parent_commandline\": \"\",\n    \"grandparent_integrity_level\": \"System\",\n    \"current_directory\": \"C:\\\\Windows\\\\system32\\\\\",\n    \"image_name\": \"C:\\\\Windows\\\\explorer.exe\",\n    \"create_time\": \"2024-01-15T08:12:43.045Z\",\n    \"ppid\": 6452,\n    \"grandparent_commandline\": \"winlogon.exe\",\n    \"log_platform_flag\": 0,\n    \"parent_commandline\": \"C:\\\\Windows\\\\system32\\\\userinit.exe\",\n    \"ancestors\": \"C:\\\\Windows\\\\System32\\\\userinit.exe|C:\\\\Windows\\\\System32\\\\winlogon.exe\",\n    \"integrity_level\": \"Medium\",\n    \"username\": \"EXAMPLE\\\\jdoe\"\n  },\n  \"alert_unique_id\": \"44c633d9-b38d-4acb-87a5-7db9bd8ab38a\",\n  \"type\": \"rtlogs\",\n  \"execution\": 0,\n  \"alert_time\": \"2024-01-15T08:13:47.621+00:00\",\n  \"image_name\": \"C:\\\\Windows\\\\explorer.exe\",\n  \"mitre_cells\": [\n    \"persistence__t1547.001\"\n  ],\n  \"@event_create_date\": \"2024-01-15T08:13:47.621Z\",\n  \"@version\": \"1\",\n  \"maturity\": \"stable\",\n  \"agent\": {\n    \"additional_info\": null,\n    \"agentid\": \"d82e1ff6-f268-42a1-9091-df230fb3c85e\",\n    \"groups\": [\n      {\n        \"id\": \"ba5a8596-5bfc-401a-a46e-1111f7d35a23\",\n        \"name\": \"EXAMPLE\"\n      }\n    ],\n    \"osversion\": \"10.0.19042\",\n    \"osproducttype\": \"Windows 10 Pro\",\n    \"domainname\": \"EXAMPLE\",\n    \"distroid\": null,\n    \"ostype\": \"windows\",\n    \"hostname\": \"PL-3049\",\n    \"version\": \"3.2.5\",\n    \"dnsdomainname\": \"example.org\",\n    \"domain\": null\n  },\n  \"@timestamp\": \"2024-01-15T08:14:59.497966Z\",\n  \"msg\": \"Detects when a file is added or modified in the startup directory\"\n}",
    "event": {
      "category": [
        "process"
      ],
      "dataset": "alert",
      "kind": "alert",
      "type": [
        "start"
      ]
    },
    "@timestamp": "2024-01-15T08:13:47.621000Z",
    "agent": {
      "id": "d82e1ff6-f268-42a1-9091-df230fb3c85e",
      "name": "harfanglab"
    },
    "file": {
      "hash": {
        "md5": "68b329da9893e34099c7d8ad5cb9c940",
        "sha1": "adc83b19e793491b1c6ea0fd8b46cd9f32e592fc",
        "sha256": "01ba4719c80b6fe911b091a7c05124b64eeece964e09c058ef8f9805daca546b"
      },
      "name": "C:\\Users\\jdoe\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\DeepL auto-start.lnk"
    },
    "harfanglab": {
      "aggregation_key": "e9cb2440a0ef88bdf8133cb7d5e41b66a903ba84df49659a05720ee1ac1caa33",
      "alert_subtype": "process",
      "alert_time": "2024-01-15T08:13:47.621+00:00",
      "alert_unique_id": "44c633d9-b38d-4acb-87a5-7db9bd8ab38a",
      "execution": 0,
      "grandparent": {
        "process": {
          "command_line": "winlogon.exe",
          "executable": "C:\\Windows\\System32\\winlogon.exe"
        }
      },
      "groups": [],
      "level": "medium",
      "status": "new",
      "threat_key": "2971"
    },
    "host": {
      "domain": "EXAMPLE",
      "hostname": "PL-3049",
      "name": "PL-3049",
      "os": {
        "full": "Windows 10 Pro",
        "version": "10.0.19042"
      }
    },
    "log": {
      "hostname": "PL-3049"
    },
    "process": {
      "command_line": "C:\\Windows\\Explorer.EXE",
      "executable": "C:\\Windows\\explorer.exe",
      "name": "explorer.exe",
      "parent": {
        "command_line": "C:\\Windows\\system32\\userinit.exe",
        "executable": "C:\\Windows\\System32\\userinit.exe"
      },
      "pe": {
        "company": "Microsoft Corporation",
        "description": "Windows Explorer",
        "file_version": "10.0.19041.1706 (WinBuild.160101.0800)",
        "original_file_name": "EXPLORER.EXE",
        "product": "Microsoft\u00ae Windows\u00ae Operating System"
      },
      "pid": 11560,
      "working_directory": "C:\\Windows\\system32\\"
    },
    "related": {
      "hash": [
        "01ba4719c80b6fe911b091a7c05124b64eeece964e09c058ef8f9805daca546b",
        "68b329da9893e34099c7d8ad5cb9c940",
        "adc83b19e793491b1c6ea0fd8b46cd9f32e592fc"
      ],
      "hosts": [
        "PL-3049"
      ],
      "user": [
        "jdoe"
      ]
    },
    "rule": {
      "category": "sigma",
      "description": "Detects when a file is added or modified in the startup directory",
      "id": "73d2962a-48ed-4fc1-937c-0a5bc688a072",
      "name": "File Added/Modified in Startup Directory"
    },
    "user": {
      "domain": "EXAMPLE",
      "name": "jdoe"
    }
  }
}