# HarfangLab/harfanglab/tests/alert_2.json
{
  "input": {
    "message": "{\n    \"type\": \"rtlogs\",\n    \"mitre_cells\": [\n        \"persistence__t1547.001\"\n    ],\n    \"level_int\": 30,\n    \"quarantine\": 4,\n    \"detection_origin\": \"agent\",\n    \"alert_type\": \"sigma\",\n    \"@timestamp\": \"2024-01-17T08:25:58.942517Z\",\n    \"details_registry\": {\n        \"details\": \"\\\"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\\\" --no-startup-window /prefetch:5\",\n        \"event_type\": \"SetValue\",\n        \"target_object\": \"HKU\\\\S-1-5-21-111111111-222222222222-3333333333-4444\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\GoogleChromeAutoLaunch_292993993B6DCE00DC4A6B0723EC700E\"\n    },\n    \"rule_id\": \"00000000-0000-0000-0000-000000000000\",\n    \"maturity\": \"stable\",\n    \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\",\n    \"alert_unique_id\": \"00000000-0000-0000-0000-000000000000\",\n    \"tenant\": \"\",\n    \"@event_create_date\": \"2024-01-17T08:19:06.071Z\",\n    \"threat_key\": 2912,\n    \"level\": \"medium\",\n    \"log_type\": \"alert\",\n    \"agent\": {\n        \"ostype\": \"windows\",\n        \"osversion\": \"10.0.19042\",\n        \"additional_info\": null,\n        \"domainname\": \"EXAMPLE\",\n        \"version\": \"3.2.9\",\n        \"domain\": null,\n        \"osproducttype\": \"Windows 10 Pro\",\n        \"dnsdomainname\": \"example.org\",\n        \"agentid\": \"00000000-0000-0000-0000-000000000000\",\n        \"distroid\": null,\n        \"hostname\": \"PL3024\"\n    },\n    \"threat_values\": [\n        \"<drive>:\\\\users\\\\<username>\\\\appdata\\\\local\\\\onelaunch\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\"\n    ],\n    \"execution\": 0,\n    \"msg\": \"Detects when an suspicious entry is added/modified in one of the autostart extensibility point (ASEP) in the registry.\\n An attacker may achieve persistence by referencing a program with a registry run key.\\n It is recommended to investigate the process that added the key as well as the target of the registry key for malicious content.\",\n    \"threat_type\": \"commandline\",\n    \"@version\": \"1\",\n    \"process\": {\n        \"parent_commandline\": \"C:\\\\Windows\\\\Explorer.EXE\",\n        \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\",\n        \"fake_parent_image\": \"\",\n        \"size\": 1802200,\n        \"pe_timestamp_int\": 1697003438,\n        \"ppid\": 3132,\n        \"grandparent_integrity_level\": \"Medium\",\n        \"fake_ppid\": 0,\n        \"dont_create_process\": true,\n        \"log_platform_flag\": 0,\n        \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\",\n        \"current_directory\": \"C:\\\\Windows\\\\system32\\\\\",\n        \"parent_integrity_level\": \"Medium\",\n        \"log_type\": \"process\",\n        \"ancestors\": \"C:\\\\Windows\\\\explorer.exe|C:\\\\Windows\\\\System32\\\\userinit.exe|C:\\\\Windows\\\\System32\\\\winlogon.exe\",\n        \"username\": \"EXAMPLE\\\\jdoe\",\n        \"fake_parent_commandline\": \"\",\n        \"grandparent_commandline\": \"C:\\\\Windows\\\\system32\\\\userinit.exe\",\n        \"parent_image\": \"C:\\\\Windows\\\\explorer.exe\",\n        \"process_name\": \"chromium.exe\",\n        \"pe_info\": {\n            \"legal_copyright\": \"Copyright 2018\",\n            \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\",\n            \"file_version\": \"118.0.0.0\",\n            \"product_version\": \"118.0.0.0\",\n            \"internal_name\": \"chromium_exe\",\n            \"file_description\": \"Chrome\",\n            \"company_name\": \"Chrome\",\n            \"original_filename\": \"chromium.exe\",\n            \"product_name\": \"Chrome\"\n        },\n        \"session\": 1,\n        \"ioc_matches\": [],\n        \"signature_info\": {\n            \"signed_authenticode\": true,\n            \"signer_info\": {\n                \"serial_number\": \"08eb9739b29536226513191ec7264032\",\n                \"issuer_name\": \"DigiCert Trusted G4 Code Signing RSA4096 SHA384 2021 CA1\",\n                \"display_name\": \"GOOGLE INC.\",\n                \"thumbprint_sha256\": \"9ef74106802ed78fc995f2b01aeaecebc1a60a7479a257f405d3520d19eaacff\",\n                \"thumbprint\": \"dcc4e5f8c45b3139dd88ce1e42a224013b81d55e\"\n            },\n            \"signed_catalog\": false,\n            \"root_info\": {\n                \"serial_number\": \"059b1b579e8e2132e23907bda777755c\",\n                \"issuer_name\": \"DigiCert Trusted Root G4\",\n                \"display_name\": \"DigiCert Trusted Root G4\",\n                \"thumbprint_sha256\": \"552f7bdcf1a7af9e6ce672017f4f12abf77240c78e761ac203d1d9d20ac89988\",\n                \"thumbprint\": \"ddfb16cd4931c973a2037d3fc83a4d7d775d05e4\"\n            }\n        },\n        \"parent_unique_id\": \"00000000-0000-0000-0000-000000000000\",\n        \"create_time\": \"2024-01-17T08:18:59.297Z\",\n        \"grandparent_image\": \"C:\\\\Windows\\\\System32\\\\userinit.exe\",\n        \"integrity_level\": \"Medium\",\n        \"commandline\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\",\n        \"pid\": 11888,\n        \"usersid\": \"S-1-5-21-111111111-222222222222-3333333333-4444\",\n        \"hashes\": {\n            \"sha1\": \"7f50d8c3cf3ec79122a876e969bdb65d939becd0\",\n            \"sha256\": \"76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4\",\n            \"md5\": \"0a4448b31ce7f83cb7691a2657f330f1\"\n        },\n        \"status\": 0,\n        \"logonid\": 3041508,\n        \"process_unique_id\": \"00000000-0000-0000-0000-000000000000\",\n        \"signed\": true\n    },\n    \"status\": \"new\",\n    \"groups\": [\n        {\n            \"name\": \"EXAMPLE\",\n            \"id\": \"00000000-0000-0000-0000-000000000000\"\n        }\n    ],\n    \"alert_time\": \"2024-01-17T08:19:06.071+00:00\",\n    \"rule_name\": \"Registry Autorun Key Added\",\n    \"alert_subtype\": \"process\",\n    \"tags\": [\n        \"attack.persistence\",\n        \"attack.t1112\",\n        \"attack.t1547.001\"\n    ]\n}"
  },
  "expected": {
    "message": "{\n    \"type\": \"rtlogs\",\n    \"mitre_cells\": [\n        \"persistence__t1547.001\"\n    ],\n    \"level_int\": 30,\n    \"quarantine\": 4,\n    \"detection_origin\": \"agent\",\n    \"alert_type\": \"sigma\",\n    \"@timestamp\": \"2024-01-17T08:25:58.942517Z\",\n    \"details_registry\": {\n        \"details\": \"\\\"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\\\" --no-startup-window /prefetch:5\",\n        \"event_type\": \"SetValue\",\n        \"target_object\": \"HKU\\\\S-1-5-21-111111111-222222222222-3333333333-4444\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\GoogleChromeAutoLaunch_292993993B6DCE00DC4A6B0723EC700E\"\n    },\n    \"rule_id\": \"00000000-0000-0000-0000-000000000000\",\n    \"maturity\": \"stable\",\n    \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\",\n    \"alert_unique_id\": \"00000000-0000-0000-0000-000000000000\",\n    \"tenant\": \"\",\n    \"@event_create_date\": \"2024-01-17T08:19:06.071Z\",\n    \"threat_key\": 2912,\n    \"level\": \"medium\",\n    \"log_type\": \"alert\",\n    \"agent\": {\n        \"ostype\": \"windows\",\n        \"osversion\": \"10.0.19042\",\n        \"additional_info\": null,\n        \"domainname\": \"EXAMPLE\",\n        \"version\": \"3.2.9\",\n        \"domain\": null,\n        \"osproducttype\": \"Windows 10 Pro\",\n        \"dnsdomainname\": \"example.org\",\n        \"agentid\": \"00000000-0000-0000-0000-000000000000\",\n        \"distroid\": null,\n        \"hostname\": \"PL3024\"\n    },\n    \"threat_values\": [\n        \"<drive>:\\\\users\\\\<username>\\\\appdata\\\\local\\\\onelaunch\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\"\n    ],\n    \"execution\": 0,\n    \"msg\": \"Detects when an suspicious entry is added/modified in one of the autostart extensibility point (ASEP) in the registry.\\n An attacker may achieve persistence by referencing a program with a registry run key.\\n It is recommended to investigate the process that added the key as well as the target of the registry key for malicious content.\",\n    \"threat_type\": \"commandline\",\n    \"@version\": \"1\",\n    \"process\": {\n        \"parent_commandline\": \"C:\\\\Windows\\\\Explorer.EXE\",\n        \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\",\n        \"fake_parent_image\": \"\",\n        \"size\": 1802200,\n        \"pe_timestamp_int\": 1697003438,\n        \"ppid\": 3132,\n        \"grandparent_integrity_level\": \"Medium\",\n        \"fake_ppid\": 0,\n        \"dont_create_process\": true,\n        \"log_platform_flag\": 0,\n        \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\",\n        \"current_directory\": \"C:\\\\Windows\\\\system32\\\\\",\n        \"parent_integrity_level\": \"Medium\",\n        \"log_type\": \"process\",\n        \"ancestors\": \"C:\\\\Windows\\\\explorer.exe|C:\\\\Windows\\\\System32\\\\userinit.exe|C:\\\\Windows\\\\System32\\\\winlogon.exe\",\n        \"username\": \"EXAMPLE\\\\jdoe\",\n        \"fake_parent_commandline\": \"\",\n        \"grandparent_commandline\": \"C:\\\\Windows\\\\system32\\\\userinit.exe\",\n        \"parent_image\": \"C:\\\\Windows\\\\explorer.exe\",\n        \"process_name\": \"chromium.exe\",\n        \"pe_info\": {\n            \"legal_copyright\": \"Copyright 2018\",\n            \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\",\n            \"file_version\": \"118.0.0.0\",\n            \"product_version\": \"118.0.0.0\",\n            \"internal_name\": \"chromium_exe\",\n            \"file_description\": \"Chrome\",\n            \"company_name\": \"Chrome\",\n            \"original_filename\": \"chromium.exe\",\n            \"product_name\": \"Chrome\"\n        },\n        \"session\": 1,\n        \"ioc_matches\": [],\n        \"signature_info\": {\n            \"signed_authenticode\": true,\n            \"signer_info\": {\n                \"serial_number\": \"08eb9739b29536226513191ec7264032\",\n                \"issuer_name\": \"DigiCert Trusted G4 Code Signing RSA4096 SHA384 2021 CA1\",\n                \"display_name\": \"GOOGLE INC.\",\n                \"thumbprint_sha256\": \"9ef74106802ed78fc995f2b01aeaecebc1a60a7479a257f405d3520d19eaacff\",\n                \"thumbprint\": \"dcc4e5f8c45b3139dd88ce1e42a224013b81d55e\"\n            },\n            \"signed_catalog\": false,\n            \"root_info\": {\n                \"serial_number\": \"059b1b579e8e2132e23907bda777755c\",\n                \"issuer_name\": \"DigiCert Trusted Root G4\",\n                \"display_name\": \"DigiCert Trusted Root G4\",\n                \"thumbprint_sha256\": \"552f7bdcf1a7af9e6ce672017f4f12abf77240c78e761ac203d1d9d20ac89988\",\n                \"thumbprint\": \"ddfb16cd4931c973a2037d3fc83a4d7d775d05e4\"\n            }\n        },\n        \"parent_unique_id\": \"00000000-0000-0000-0000-000000000000\",\n        \"create_time\": \"2024-01-17T08:18:59.297Z\",\n        \"grandparent_image\": \"C:\\\\Windows\\\\System32\\\\userinit.exe\",\n        \"integrity_level\": \"Medium\",\n        \"commandline\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\",\n        \"pid\": 11888,\n        \"usersid\": \"S-1-5-21-111111111-222222222222-3333333333-4444\",\n        \"hashes\": {\n            \"sha1\": \"7f50d8c3cf3ec79122a876e969bdb65d939becd0\",\n            \"sha256\": \"76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4\",\n            \"md5\": \"0a4448b31ce7f83cb7691a2657f330f1\"\n        },\n        \"status\": 0,\n        \"logonid\": 3041508,\n        \"process_unique_id\": \"00000000-0000-0000-0000-000000000000\",\n        \"signed\": true\n    },\n    \"status\": \"new\",\n    \"groups\": [\n        {\n            \"name\": \"EXAMPLE\",\n            \"id\": \"00000000-0000-0000-0000-000000000000\"\n        }\n    ],\n    \"alert_time\": \"2024-01-17T08:19:06.071+00:00\",\n    \"rule_name\": \"Registry Autorun Key Added\",\n    \"alert_subtype\": \"process\",\n    \"tags\": [\n        \"attack.persistence\",\n        \"attack.t1112\",\n        \"attack.t1547.001\"\n    ]\n}",
    "event": {
      "category": [
        "process"
      ],
      "dataset": "alert",
      "kind": "alert",
      "type": [
        "start"
      ]
    },
    "@timestamp": "2024-01-17T08:19:06.071000Z",
    "agent": {
      "id": "00000000-0000-0000-0000-000000000000",
      "name": "harfanglab"
    },
    "file": {
      "hash": {
        "md5": "0a4448b31ce7f83cb7691a2657f330f1",
        "sha1": "7f50d8c3cf3ec79122a876e969bdb65d939becd0",
        "sha256": "76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4"
      }
    },
    "harfanglab": {
      "alert_subtype": "process",
      "alert_time": "2024-01-17T08:19:06.071+00:00",
      "alert_unique_id": "00000000-0000-0000-0000-000000000000",
      "execution": 0,
      "grandparent": {
        "process": {
          "command_line": "C:\\Windows\\system32\\userinit.exe",
          "executable": "C:\\Windows\\System32\\userinit.exe"
        }
      },
      "groups": [
        "{\"id\":\"00000000-0000-0000-0000-000000000000\",\"name\":\"EXAMPLE\"}"
      ],
      "level": "medium",
      "status": "new",
      "threat_key": "2912"
    },
    "host": {
      "domain": "EXAMPLE",
      "hostname": "PL3024",
      "name": "PL3024",
      "os": {
        "full": "Windows 10 Pro",
        "version": "10.0.19042"
      }
    },
    "log": {
      "hostname": "PL3024"
    },
    "process": {
      "command_line": "C:\\Users\\jdoe\\AppData\\Local\\Chrome\\5.25.0\\chromium\\chromium.exe --no-startup-window /prefetch:5",
      "executable": "C:\\Users\\jdoe\\AppData\\Local\\Chrome\\5.25.0\\chromium\\chromium.exe",
      "name": "chromium.exe",
      "parent": {
        "command_line": "C:\\Windows\\Explorer.EXE",
        "executable": "C:\\Windows\\explorer.exe"
      },
      "pe": {
        "company": "Chrome",
        "description": "Chrome",
        "file_version": "118.0.0.0",
        "original_file_name": "chromium.exe",
        "product": "Chrome"
      },
      "pid": 11888,
      "working_directory": "C:\\Windows\\system32\\"
    },
    "registry": {
      "hive": "HKU",
      "key": "S-1-5-21-111111111-222222222222-3333333333-4444\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run",
      "path": "HKU\\S-1-5-21-111111111-222222222222-3333333333-4444\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\GoogleChromeAutoLaunch_292993993B6DCE00DC4A6B0723EC700E",
      "value": "\"C:\\Users\\jdoe\\AppData\\Local\\Chrome\\5.25.0\\chromium\\chromium.exe\" --no-startup-window /prefetch:5"
    },
    "related": {
      "hash": [
        "0a4448b31ce7f83cb7691a2657f330f1",
        "76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4",
        "7f50d8c3cf3ec79122a876e969bdb65d939becd0"
      ],
      "hosts": [
        "PL3024"
      ],
      "user": [
        "jdoe"
      ]
    },
    "rule": {
      "category": "sigma",
      "description": "Detects when an suspicious entry is added/modified in one of the autostart extensibility point (ASEP) in the registry.\n An attacker may achieve persistence by referencing a program with a registry run key.\n It is recommended to investigate the process that added the key as well as the target of the registry key for malicious content.",
      "id": "00000000-0000-0000-0000-000000000000",
      "name": "Registry Autorun Key Added"
    },
    "user": {
      "domain": "EXAMPLE",
      "name": "jdoe",
      "roles": "EXAMPLE"
    }
  }
}