# HarfangLab/harfanglab/tests/alert_2.json
{
  "input": {
    "message": "{\"type\": \"rtlogs\", \"mitre_cells\": [\"persistence__t1547.001\"], \"level_int\": 30, \"quarantine\": 4, \"detection_origin\": \"agent\", \"alert_type\": \"sigma\", \"@timestamp\": \"2024-01-17T08:25:58.942517Z\", \"details_registry\": {\"details\": \"\\\"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\\\" --no-startup-window /prefetch:5\", \"event_type\": \"SetValue\", \"target_object\": \"HKU\\\\S-1-5-21-111111111-222222222222-3333333333-4444\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\GoogleChromeAutoLaunch_292993993B6DCE00DC4A6B0723EC700E\"}, \"rule_id\": \"00000000-0000-0000-0000-000000000000\", \"maturity\": \"stable\", \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\", \"alert_unique_id\": \"00000000-0000-0000-0000-000000000000\", \"tenant\": \"\", \"@event_create_date\": \"2024-01-17T08:19:06.071Z\", \"threat_key\": 2912, \"level\": \"medium\", \"log_type\": \"alert\", \"agent\": {\"ostype\": \"windows\", \"osversion\": \"10.0.19042\", \"additional_info\": null, \"domainname\": \"EXAMPLE\", \"version\": \"3.2.9\", \"domain\": null, \"osproducttype\": \"Windows 10 Pro\", \"dnsdomainname\": \"example.org\", \"agentid\": \"00000000-0000-0000-0000-000000000000\", \"distroid\": null, \"hostname\": \"PL3024\"}, \"threat_values\": [\"<drive>:\\\\users\\\\<username>\\\\appdata\\\\local\\\\onelaunch\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\"], \"execution\": 0, \"msg\": \"Detects when an suspicious entry is added/modified in one of the autostart extensibility point (ASEP) in the registry.\\n An attacker may achieve persistence by referencing a program with a registry run key.\\n It is recommended to investigate the process that added the key as well as the target of the registry key for malicious content.\", \"threat_type\": \"commandline\", \"@version\": \"1\", \"process\": {\"parent_commandline\": \"C:\\\\Windows\\\\Explorer.EXE\", \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\", \"fake_parent_image\": \"\", \"size\": 1802200, \"pe_timestamp_int\": 1697003438, \"ppid\": 3132, \"grandparent_integrity_level\": \"Medium\", \"fake_ppid\": 0, \"dont_create_process\": true, \"log_platform_flag\": 0, \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\", \"current_directory\": \"C:\\\\Windows\\\\system32\\\\\", \"parent_integrity_level\": \"Medium\", \"log_type\": \"process\", \"ancestors\": \"C:\\\\Windows\\\\explorer.exe|C:\\\\Windows\\\\System32\\\\userinit.exe|C:\\\\Windows\\\\System32\\\\winlogon.exe\", \"username\": \"EXAMPLE\\\\jdoe\", \"fake_parent_commandline\": \"\", \"grandparent_commandline\": \"C:\\\\Windows\\\\system32\\\\userinit.exe\", \"parent_image\": \"C:\\\\Windows\\\\explorer.exe\", \"process_name\": \"chromium.exe\", \"pe_info\": {\"legal_copyright\": \"Copyright 2018\", \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\", \"file_version\": \"118.0.0.0\", \"product_version\": \"118.0.0.0\", \"internal_name\": \"chromium_exe\", \"file_description\": \"Chrome\", \"company_name\": \"Chrome\", \"original_filename\": \"chromium.exe\", \"product_name\": \"Chrome\"}, \"session\": 1, \"ioc_matches\": [], \"signature_info\": {\"signed_authenticode\": true, \"signer_info\": {\"serial_number\": \"08eb9739b29536226513191ec7264032\", \"issuer_name\": \"DigiCert Trusted G4 Code Signing RSA4096 SHA384 2021 CA1\", \"display_name\": \"GOOGLE INC.\", \"thumbprint_sha256\": \"9ef74106802ed78fc995f2b01aeaecebc1a60a7479a257f405d3520d19eaacff\", \"thumbprint\": \"dcc4e5f8c45b3139dd88ce1e42a224013b81d55e\"}, \"signed_catalog\": false, \"root_info\": {\"serial_number\": \"059b1b579e8e2132e23907bda777755c\", \"issuer_name\": \"DigiCert Trusted Root G4\", \"display_name\": \"DigiCert Trusted Root G4\", \"thumbprint_sha256\": \"552f7bdcf1a7af9e6ce672017f4f12abf77240c78e761ac203d1d9d20ac89988\", \"thumbprint\": \"ddfb16cd4931c973a2037d3fc83a4d7d775d05e4\"}}, \"parent_unique_id\": \"00000000-0000-0000-0000-000000000000\", \"create_time\": \"2024-01-17T08:18:59.297Z\", \"grandparent_image\": \"C:\\\\Windows\\\\System32\\\\userinit.exe\", \"integrity_level\": \"Medium\", \"commandline\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\", \"pid\": 11888, \"usersid\": \"S-1-5-21-111111111-222222222222-3333333333-4444\", \"hashes\": {\"sha1\": \"7f50d8c3cf3ec79122a876e969bdb65d939becd0\", \"sha256\": \"76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4\", \"md5\": \"0a4448b31ce7f83cb7691a2657f330f1\"}, \"status\": 0, \"logonid\": 3041508, \"process_unique_id\": \"00000000-0000-0000-0000-000000000000\", \"signed\": true}, \"status\": \"new\", \"groups\": [{\"name\": \"EXAMPLE\", \"id\": \"00000000-0000-0000-0000-000000000000\"}], \"alert_time\": \"2024-01-17T08:19:06.071+00:00\", \"rule_name\": \"Registry Autorun Key Added\", \"alert_subtype\": \"process\", \"tags\": [\"attack.persistence\", \"attack.t1112\", \"attack.t1547.001\"]}",
    "sekoiaio": {
      "intake": {
        "dialect": "harfanglab",
        "dialect_uuid": "3c7057d3-4689-4fae-8033-6f1f887a70f2"
      }
    }
  },
  "expected": {
    "message": "{\"type\": \"rtlogs\", \"mitre_cells\": [\"persistence__t1547.001\"], \"level_int\": 30, \"quarantine\": 4, \"detection_origin\": \"agent\", \"alert_type\": \"sigma\", \"@timestamp\": \"2024-01-17T08:25:58.942517Z\", \"details_registry\": {\"details\": \"\\\"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\\\" --no-startup-window /prefetch:5\", \"event_type\": \"SetValue\", \"target_object\": \"HKU\\\\S-1-5-21-111111111-222222222222-3333333333-4444\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\GoogleChromeAutoLaunch_292993993B6DCE00DC4A6B0723EC700E\"}, \"rule_id\": \"00000000-0000-0000-0000-000000000000\", \"maturity\": \"stable\", \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\", \"alert_unique_id\": \"00000000-0000-0000-0000-000000000000\", \"tenant\": \"\", \"@event_create_date\": \"2024-01-17T08:19:06.071Z\", \"threat_key\": 2912, \"level\": \"medium\", \"log_type\": \"alert\", \"agent\": {\"ostype\": \"windows\", \"osversion\": \"10.0.19042\", \"additional_info\": null, \"domainname\": \"EXAMPLE\", \"version\": \"3.2.9\", \"domain\": null, \"osproducttype\": \"Windows 10 Pro\", \"dnsdomainname\": \"example.org\", \"agentid\": \"00000000-0000-0000-0000-000000000000\", \"distroid\": null, \"hostname\": \"PL3024\"}, \"threat_values\": [\"<drive>:\\\\users\\\\<username>\\\\appdata\\\\local\\\\onelaunch\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\"], \"execution\": 0, \"msg\": \"Detects when an suspicious entry is added/modified in one of the autostart extensibility point (ASEP) in the registry.\\n An attacker may achieve persistence by referencing a program with a registry run key.\\n It is recommended to investigate the process that added the key as well as the target of the registry key for malicious content.\", \"threat_type\": \"commandline\", \"@version\": \"1\", \"process\": {\"parent_commandline\": \"C:\\\\Windows\\\\Explorer.EXE\", \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\", \"fake_parent_image\": \"\", \"size\": 1802200, \"pe_timestamp_int\": 1697003438, \"ppid\": 3132, \"grandparent_integrity_level\": \"Medium\", \"fake_ppid\": 0, \"dont_create_process\": true, \"log_platform_flag\": 0, \"image_name\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe\", \"current_directory\": \"C:\\\\Windows\\\\system32\\\\\", \"parent_integrity_level\": \"Medium\", \"log_type\": \"process\", \"ancestors\": \"C:\\\\Windows\\\\explorer.exe|C:\\\\Windows\\\\System32\\\\userinit.exe|C:\\\\Windows\\\\System32\\\\winlogon.exe\", \"username\": \"EXAMPLE\\\\jdoe\", \"fake_parent_commandline\": \"\", \"grandparent_commandline\": \"C:\\\\Windows\\\\system32\\\\userinit.exe\", \"parent_image\": \"C:\\\\Windows\\\\explorer.exe\", \"process_name\": \"chromium.exe\", \"pe_info\": {\"legal_copyright\": \"Copyright 2018\", \"pe_timestamp\": \"2023-10-11T05:50:38.000Z\", \"file_version\": \"118.0.0.0\", \"product_version\": \"118.0.0.0\", \"internal_name\": \"chromium_exe\", \"file_description\": \"Chrome\", \"company_name\": \"Chrome\", \"original_filename\": \"chromium.exe\", \"product_name\": \"Chrome\"}, \"session\": 1, \"ioc_matches\": [], \"signature_info\": {\"signed_authenticode\": true, \"signer_info\": {\"serial_number\": \"08eb9739b29536226513191ec7264032\", \"issuer_name\": \"DigiCert Trusted G4 Code Signing RSA4096 SHA384 2021 CA1\", \"display_name\": \"GOOGLE INC.\", \"thumbprint_sha256\": \"9ef74106802ed78fc995f2b01aeaecebc1a60a7479a257f405d3520d19eaacff\", \"thumbprint\": \"dcc4e5f8c45b3139dd88ce1e42a224013b81d55e\"}, \"signed_catalog\": false, \"root_info\": {\"serial_number\": \"059b1b579e8e2132e23907bda777755c\", \"issuer_name\": \"DigiCert Trusted Root G4\", \"display_name\": \"DigiCert Trusted Root G4\", \"thumbprint_sha256\": \"552f7bdcf1a7af9e6ce672017f4f12abf77240c78e761ac203d1d9d20ac89988\", \"thumbprint\": \"ddfb16cd4931c973a2037d3fc83a4d7d775d05e4\"}}, \"parent_unique_id\": \"00000000-0000-0000-0000-000000000000\", \"create_time\": \"2024-01-17T08:18:59.297Z\", \"grandparent_image\": \"C:\\\\Windows\\\\System32\\\\userinit.exe\", \"integrity_level\": \"Medium\", \"commandline\": \"C:\\\\Users\\\\jdoe\\\\AppData\\\\Local\\\\Chrome\\\\5.25.0\\\\chromium\\\\chromium.exe --no-startup-window /prefetch:5\", \"pid\": 11888, \"usersid\": \"S-1-5-21-111111111-222222222222-3333333333-4444\", \"hashes\": {\"sha1\": \"7f50d8c3cf3ec79122a876e969bdb65d939becd0\", \"sha256\": \"76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4\", \"md5\": \"0a4448b31ce7f83cb7691a2657f330f1\"}, \"status\": 0, \"logonid\": 3041508, \"process_unique_id\": \"00000000-0000-0000-0000-000000000000\", \"signed\": true}, \"status\": \"new\", \"groups\": [{\"name\": \"EXAMPLE\", \"id\": \"00000000-0000-0000-0000-000000000000\"}], \"alert_time\": \"2024-01-17T08:19:06.071+00:00\", \"rule_name\": \"Registry Autorun Key Added\", \"alert_subtype\": \"process\", \"tags\": [\"attack.persistence\", \"attack.t1112\", \"attack.t1547.001\"]}",
    "event": {
      "category": [
        "process"
      ],
      "dataset": "alert",
      "kind": "alert",
      "type": [
        "start"
      ]
    },
    "@timestamp": "2024-01-17T08:19:06.071000Z",
    "agent": {
      "id": "00000000-0000-0000-0000-000000000000",
      "name": "harfanglab"
    },
    "file": {
      "hash": {
        "md5": "0a4448b31ce7f83cb7691a2657f330f1",
        "sha1": "7f50d8c3cf3ec79122a876e969bdb65d939becd0",
        "sha256": "76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4"
      }
    },
    "harfanglab": {
      "alert_subtype": "process",
      "alert_time": "2024-01-17T08:19:06.071+00:00",
      "alert_unique_id": "00000000-0000-0000-0000-000000000000",
      "execution": 0,
      "groups": [
        "{\"id\": \"00000000-0000-0000-0000-000000000000\", \"name\": \"EXAMPLE\"}"
      ],
      "level": "medium",
      "status": "new"
    },
    "host": {
      "domain": "EXAMPLE",
      "hostname": "PL3024",
      "name": "PL3024",
      "os": {
        "full": "Windows 10 Pro",
        "version": "10.0.19042"
      }
    },
    "log": {
      "hostname": "PL3024"
    },
    "process": {
      "command_line": "C:\\Users\\jdoe\\AppData\\Local\\Chrome\\5.25.0\\chromium\\chromium.exe --no-startup-window /prefetch:5",
      "executable": "C:\\Users\\jdoe\\AppData\\Local\\Chrome\\5.25.0\\chromium\\chromium.exe",
      "name": "chromium.exe",
      "parent": {
        "command_line": "C:\\Windows\\Explorer.EXE",
        "executable": "C:\\Windows\\explorer.exe"
      },
      "pe": {
        "company": "Chrome",
        "description": "Chrome",
        "file_version": "118.0.0.0",
        "original_file_name": "chromium.exe",
        "product": "Chrome"
      },
      "pid": 11888,
      "working_directory": "C:\\Windows\\system32\\"
    },
    "registry": {
      "hive": "HKU",
      "key": "S-1-5-21-111111111-222222222222-3333333333-4444\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run",
      "path": "HKU\\S-1-5-21-111111111-222222222222-3333333333-4444\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\GoogleChromeAutoLaunch_292993993B6DCE00DC4A6B0723EC700E",
      "value": "\"C:\\Users\\jdoe\\AppData\\Local\\Chrome\\5.25.0\\chromium\\chromium.exe\" --no-startup-window /prefetch:5"
    },
    "related": {
      "hash": [
        "0a4448b31ce7f83cb7691a2657f330f1",
        "76eac7b5f53e0d58a98d5a6ddf9c97e19d1462ef65c0035d7798f89988b15ab4",
        "7f50d8c3cf3ec79122a876e969bdb65d939becd0"
      ],
      "hosts": [
        "PL3024"
      ],
      "user": [
        "EXAMPLE\\jdoe"
      ]
    },
    "rule": {
      "category": "sigma",
      "description": "Detects when an suspicious entry is added/modified in one of the autostart extensibility point (ASEP) in the registry.\n An attacker may achieve persistence by referencing a program with a registry run key.\n It is recommended to investigate the process that added the key as well as the target of the registry key for malicious content.",
      "id": "00000000-0000-0000-0000-000000000000",
      "name": "Registry Autorun Key Added"
    },
    "user": {
      "name": "EXAMPLE\\jdoe",
      "roles": "EXAMPLE"
    }
  }
}